# deploy/systemd/crabtrap-lite.service
[Unit]
Description=CrabTrap Action Gate
After=network-online.target crabtrap-llm-reviewer.service
Wants=network-online.target crabtrap-llm-reviewer.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/crabtrap-lite
EnvironmentFile=/root/.config/crabtrap-lite.env
ExecStart=/usr/bin/node /root/projects/crabtrap-lite/src/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
